# cs655.service
[Unit]
Description=Run Forever

[Service]
ExecStart=/users/fspaeh/cs655-pw-cracker/src/start.sh

[Install]
WantedBy=multi-user.target
